# .config/systemd/user/zfs-auto-snapshot@hourly.timer
[Unit]
Description=ZFS %I auto snapshot timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
