# .config/systemd/user/zfs-auto-snapshot@weekly.timer
[Unit]
Description=ZFS %I auto snapshot timer

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
